# systemd/printstream.service
[Unit]
Description=PrintStream
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/printstream
EnvironmentFile=/opt/printstream/.env
ExecStart=/opt/printstream/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8088
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
